# lib/systemd/system/mq200embed.service
[Unit]
Description=BCS mq200embed

Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
ExecStart=/OLMA/mq200e_em
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
